# config/default/rootfs/etc/rc.local
#!/bin/sh -l

exec 0<>/dev/tty9
exec 1>&0
exec 2>&0

printf '\033%%G' >/dev/tty9

chvt 9
env - LC_ALL=C.UTF-8 TERM=linux \
  PATH=/sbin:/usr/sbin:/bin/:/usr/bin/ \
  setsid -c console-keyboard-multiplexer -w : -- /root/first_boot_setup/setup.sh

# Remove one time script
set +x
rm /etc/rc.local
cat >/etc/rc.local <<EOF
#!/bin/sh
exit 0
EOF
exit 0
