# rootfs_custom_files/etc/rc.local
#!/bin/sh

/root/first_boot_setup.sh

# Remove one time script
set +x
rm /etc/rc.local
cat >/etc/rc.local <<EOF
#!/bin/sh
exit 0
EOF
exit 0
