# config/devuan/rootfs/etc/init.d/wys
#!/lib/init/init-d-script
### BEGIN INIT INFO
# Provides:          wys
# Required-Start:    $local_fs pulseaudio dbus
# Required-Stop:     $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: wys daemon
# Description:       wys daemon
### END INIT INFO

DAEMON=/usr/bin/wys
PIDFILE=/var/run/wys.pid
START_ARGS="-c wys --background -m"
